# root/etc/s6-overlay/s6-rc.d/ez-share-cron/up
/etc/s6-overlay/s6-rc.d/ez-share-cron/run
